# wip_ramdisk/etc/inittab
# if this is the only init - this will print the banner once - and without a tty - you will 
::once:/bin/cat /extra/banner

# if done like this and you exit the shell - you will not get another shell again:
# ::sysinit:/bin/sh 

# if done like this and you exit the shell - you will get another shell:
::respawn:/bin/sh


